# stasher/packaging/freebsd9/databases/stasher/files/stasherd.in
#! /bin/sh
#
#

# PROVIDE: stasherd
# REQUIRE: httportmap
# KEYWORD: shutdown

#
# Add the following lines to /etc/rc.conf to enable stasherd:
#
#stasherd_enable="YES"

. /etc/rc.subr

name="stasherd"
rcvar=`set_rcvar`

load_rc_config $name

start_cmd="stasherd_start"
stop_cmd="stasherd_stop"

prefix="%%PREFIX%%"

stasherd_start() {

    if test ! -f $prefix/sbin/httpd
    then
	stasherd_args="$stasherd_args --http"
    fi
    $prefix/share/stasher/sysinit.rc $prefix/share/stasher/stasher.sysinit start
}

stasherd_stop() {

    $prefix/share/stasher/sysinit.rc $prefix/share/stasher/stasher.sysinit stop
}
run_rc_command "$1"
